# packages/tradfi/deploy/ogg-tradfi-hourly-restart.service
[Unit]
Description=Hourly restart of ogg-tradfi.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart ogg-tradfi.service
