# packages/tradfi/deploy/ogg-tradfi-hourly-restart.timer
[Unit]
Description=Restart ogg-tradfi.service every hour

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
